# packages/bsp/common/etc/NetworkManager/dispatcher.d/80-update-htop-and-offload-tx
#!/bin/bash
#
# adjust htop settings for all normal users and root
# Please fix this inside htop and DO NOT write into users config files every bootup!
